# packaging/systemd/bor-server.service
[Unit]
Description=Bor Policy Management Server
Documentation=https://github.com/VuteTech/Bor
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=bor
Group=bor

# Primary configuration is /etc/bor/server.yaml (read directly by the binary).
# The env file is optional — useful for secrets injection from Vault, k8s, etc.
# Remove the leading '-' to make it mandatory.
EnvironmentFile=-/etc/bor/server.env

ExecStart=/usr/bin/bor-server
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bor-server

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/bor/pki /var/lib/bor/acme

[Install]
WantedBy=multi-user.target
